# provisioning/roles/supervisor/templates/gunicorn.conf
[program:gunicorn]
command=/home/{{ user }}/gunicorn.sh
autostart=true
autorestart=true
stderr_logfile=/home/{{ user }}/logs/gunicorn.err.log
stdout_logfile=/home/{{ user }}/logs/gunicorn.out.log