# provisioning/roles/supervisor/templates/long_script.conf
[program:long_script]
command=/home/{{ user }}/long.sh
autostart=true
autorestart=true
stderr_logfile=/home/{{ user }}/logs/long.err.log
stdout_logfile=/home/{{ user }}/logs/long.out.log